# TestFolder/systemd/ClassNotice.timer
[Unit]
Description=Run Test.service at 08:40

[Timer]
OnCalendar=08:40:00
Unit=Test.service

[Install]
WantedBy=timers.target
